# SpiralLinux-Budgie/config/includes.chroot/etc/skel/.config/autostart/nemo-desktop.desktop
[Desktop Entry]
Type=Application
Name=Nemo Desktop
Description=Nemo desktop icons
Exec=/usr/bin/nemo-desktop
