# SpiralLinux-LXQt/config/includes.chroot/etc/skel/.config/autostart/xsettingsd.desktop
[Desktop Entry]
Exec=xsettingsd
Name=Xsettingsd (for GTK font rendering)
OnlyShowIn=LXQt;
Type=Application
Version=1.0
