# SpiralLinux-Plasma/config/includes.chroot/etc/skel/.config/autostart/org.kde.discover.notifier.desktop
[Desktop Entry]
Name=Disable Discover Notifier applet
Exec=/usr/lib/x86_64-linux-gnu/libexec/DiscoverNotifier -x
Icon=system-software-update
Type=Application
NoDisplay=true
OnlyShowIn=KDE
