# SpiralLinux-Gnome/config/includes.chroot/etc/skel/.config/autostart/trust-desktop-icons.desktop
[Desktop Entry]
Exec=/usr/bin/trust-desktop-icons.sh
Name=Trust Desktop icons
Type=Application
Hidden=false
X-GNOME-Autostart-enabled=true
